# config/system/com.cbi.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
WATCHDOG - LaunchDaemon

This daemon runs the watchdog script every hour to monitor pipeline health.
It checks for missed runs, failures, and triggers recovery if needed.

Installation:
  sudo cp config/system/com.cbi.watchdog.plist /Library/LaunchDaemons/
  sudo launchctl load /Library/LaunchDaemons/com.cbi.watchdog.plist

Uninstallation:
  sudo launchctl unload /Library/LaunchDaemons/com.cbi.watchdog.plist
  sudo rm /Library/LaunchDaemons/com.cbi.watchdog.plist

Check status:
  sudo launchctl list | grep com.cbi.watchdog
  
View logs:
  cat ~/Documents/GitHub/CBI-V14/Logs/nightly/watchdog.log
-->

<plist version="1.0">
<dict>
   <!-- Service identifier -->
   <key>Label</key>
   <string>com.cbi.watchdog</string>

   <!-- Command to run -->
   <key>ProgramArguments</key>
   <array>
       <string>/usr/local/bin/python3</string>
       <string>/Users/kirkmusick/Documents/GitHub/CBI-V14/scripts/watchdog.py</string>
   </array>

   <!-- Schedule: Run every hour -->
   <key>StartInterval</key>
   <integer>3600</integer>

   <!-- Logging -->
   <key>StandardOutPath</key>
   <string>/Users/kirkmusick/Documents/GitHub/CBI-V14/Logs/nightly/watchdog_stdout.log</string>

   <key>StandardErrorPath</key>
   <string>/Users/kirkmusick/Documents/GitHub/CBI-V14/Logs/nightly/watchdog_stderr.log</string>

   <!-- Run on load (start immediately) -->
   <key>RunAtLoad</key>
   <true/>

   <!-- Don't keep alive - run on schedule -->
   <key>KeepAlive</key>
   <false/>

   <!-- Working directory -->
   <key>WorkingDirectory</key>
   <string>/Users/kirkmusick/Documents/GitHub/CBI-V14</string>

   <!-- Environment variables -->
   <key>EnvironmentVariables</key>
   <dict>
     <key>PROJECT</key>
     <string>cbi-v14</string>
     <key>PATH</key>
     <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
   </dict>

</dict>
</plist>







